# roles/techlog/files/techlog-archive.timer
[Unit]
Description=Timer for techlog-archive

[Timer]
OnCalendar=00/1:10:00
Unit=techlog-archive.service

[Install]
WantedBy=multi-user.target
